# led-matrix-server.service
[Unit]
Description=led-matrix-server (HTTP control server + auto-cycler)
# The stock ticker mode needs the network to fetch quotes.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# The binary reads all settings from led-matrix-server.conf next to it and takes
# no arguments. cwd is set so modes find fonts/ (and the conf) via relative paths.
WorkingDirectory=/root/led-matrix-server
ExecStart=/root/led-matrix-server/led-matrix-server
# GPIO/DMA needs root; the app keeps root (privilege dropping is disabled for
# the static binary) and asks for slightly elevated scheduling for a steady
# refresh rate.
User=root
Nice=-10
Restart=on-failure
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
